# roles/prometheus/templates/etc/systemd/system/vault.service.d/override.conf
# {{ ansible_managed }}

[Service]
ExecStart=
ExecStart=/usr/bin/vault proxy -config=/etc/vault.d/vault.hcl